# init/cabled
#!/sbin/runscript

description="Starts cables communication daemon."

cuser=cable
cableflag=/srv/www/cable

daemon=/usr/libexec/cable/cabled
pidfile=/var/run/cabled.pid


start() {
    if [ ! -e ${cableflag} ]; then
        # SSD -u sets up USER and HOME, and sanitizes environment
        ebegin "Starting cabled"
        start-stop-daemon -S -q -u ${cuser} -p ${pidfile} -bm -x ${daemon} --
        eend $?
    else
        ewarn "Skipping cabled startup"
    fi
}

stop() {
    if [ ! -e ${cableflag} ]; then
        ebegin "Stopping cabled"
        start-stop-daemon -K -q -u ${cuser} -p ${pidfile}
        eend $?
    fi
}
