# conf/cabled
#!/sbin/runscript

description="Cables communication daemon"

command=/usr/libexec/cable/cabled
command_background="true"
start_stop_daemon_args="-u cable"
pidfile=/var/run/cabled.pid
